# systemd/system/upgrade-ytdl.timer
[Unit]
Description=Upgrade YoutubeDl weekly

[Timer]
OnCalendar=Fri 14:00
Persistent=true

[Install]
WantedBy=timers.target
